# packages/ytl-linux-digabi2-examnet/docker.service.d/wait-for-listen-ip.conf
[Unit]
After=ytl-linux-digabi2-wait-for-listen-ip.service
Wants=ytl-linux-digabi2-wait-for-listen-ip.service
