# packages/ytl-linux-digabi2-examnet/ytl-linux-digabi2-examnet-discovery.timer
[Unit]
Description=Periodically scan for aliases in Digabi 2 exam network

[Timer]
AccuracySec=1seconds
OnUnitInactiveSec=1minutes

[Install]
WantedBy=multi-user.target
