# packages/ytl-linux-digabi2-examnet/ytl-linux-digabi2-examnet.service
[Unit]
Description=Digabi 2 examnet background services

[Service]
ExecStart=/usr/local/sbin/digabi2-examnet-bouncer # (...)
WorkingDirectory=/tmp
StandardOutput=append:/dev/null
StandardError=append:/dev/null
Restart=always

[Install]
WantedBy=multi-user.target
