# packages/ytl-linux-update-checker/ytl-linux-update-check.service
[Unit]
Description=Automatic update checker for YTL Linux

[Service]
ExecStart=/usr/local/sbin/ytl-linux-update-checker
StandardOutput=journal
StandardError=journal
