# packages/ytl-linux-update-checker/ytl-linux-update-check.timer
[Unit]
Description=Periodically check for available updates to YTL Linux

[Timer]
OnUnitInactiveSec=1hour

[Install]
WantedBy=multi-user.target
