# packages/ytl-linux-update-checker/ytl-linux-update-check.timer
[Unit]
Description=Periodically check for available updates to YTL Linux

[Timer]
OnStartupSec=90seconds
OnUnitInactiveSec=1hour

[Install]
WantedBy=timers.target
