# packages/ytl-linux-tasks/ytl-tasks-maintenance.service
[Unit]
Description=YTL Linux maintenance tasks

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ytl-tasks maintenance

[Install]
WantedBy=multi-user.target
